# docker/php/crontab
* * * * * /usr/local/bin/php /var/www/project/bin/console simulations-cleanup
